# printers/PM-241-BT/label-printer.service
[Unit]
Description=Phomemo Label Printer Service
After=network.target cups.service
Wants=cups.service

[Service]
Type=simple
User=nthmost
Group=nthmost
WorkingDirectory=/home/nthmost/label_printer
Environment="PATH=/home/nthmost/label_venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/nthmost/label_printer/.env
ExecStart=/home/nthmost/label_venv/bin/python3 /home/nthmost/label_printer/label_service.py --host 0.0.0.0 --port 8001

# Restart behavior
Restart=on-failure
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=300

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=label-printer

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/tmp/label_images

[Install]
WantedBy=multi-user.target
